# sv/supervisor.conf.example
[program:dcts]
directory=/home/dcts
command=sh check.sh
autostart=true
autorestart=true
stderr_logfile=/home/dcts/err.log
stdout_logfile=/home/dcts/out.log
stopsignal=QUIT