# sv/supervisor.conf.example
[program:dcts]
directory=/home/dcts/sv
command=bash check.sh
autostart=true
autorestart=true
stderr_logfile=/home/dcts/sv/err.log
stdout_logfile=/home/dcts/sv/out.log
stopsignal=QUIT